# com.readable.tts.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.readable.tts</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/jesper/Projects/Dev_projects/Tools/Readable/.venv/bin/python</string>
        <string>-m</string>
        <string>readable.app_optimized</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/jesper/Projects/Dev_projects/Tools/Readable</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/jesper/.readable/logs/launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jesper/.readable/logs/launchd.error.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
